# Anim_clock/rpi2dmd.service
[Unit]
Description=AnimClock

[Service]
ExecStart=/opt/AnimClock/animclock.sh

[Install]
WantedBy=multi-user.target

