# serum/serum.service
[Unit]
Description=spi-serum

[Service]
ExecStart=/usr/local/bin/spi-serum /boot/firmware/ serum

[Install]
WantedBy=multi-user.target

