# serum/serum.service
[Unit]
Description=spi-serum

[Service]
ExecStart=/opt/serum/spi-serum /media/data/ serum

[Install]
WantedBy=multi-user.target

